# supervisor/conf.d/voicecli_tts.conf
[program:voicecli_tts]
command=%(ENV_HOME)s/.local/bin/voicecli serve --engine qwen-fast
directory=%(ENV_HOME)s/projects/voiceCLI
environment=HOME="%(ENV_HOME)s",PATH="%(ENV_HOME)s/.local/bin:%(ENV_HOME)s/projects/voiceCLI/.venv/bin:/usr/local/bin:/usr/bin:/bin"
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/projects/voiceCLI/supervisor/logs/voicecli_tts.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=%(ENV_HOME)s/projects/voiceCLI/supervisor/logs/voicecli_tts_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
